# task-18-backup/ansible/roles/client/templates/borg-backup.timer.j2
[Unit]
Description=Borg Backup Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec={{ backup_timer_interval }}
Persistent=true

[Install]
WantedBy=timers.target
